# supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=/usr/bin/mysqld_safe --datadir=/data/mysql
autostart=true
autorestart=true
stdout_logfile=/data/logs/mysql.log
stderr_logfile=/data/logs/mysql.err

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
stdout_logfile=/data/logs/nginx.log
stderr_logfile=/data/logs/nginx.err

[program:php-fpm]
command=/usr/sbin/php-fpm
autostart=true
autorestart=true
stdout_logfile=/data/logs/php-fpm.log
stderr_logfile=/data/logs/php-fpm.err
